# scripts/docker-fail2ban-discord.cron.example
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

# fail2ban の BAN / 解除を Discord へ毎時通知します。
INFRA_FAIL2BAN_DIR=/path/to/infra-fail2ban
0 * * * * root ${INFRA_FAIL2BAN_DIR}/scripts/docker-fail2ban-discord.sh
